# cron-tasks/systemd/yoda-cron@.service.template
[Unit]
Description=Yoda cron task: %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{INSTALL_DIR}}/workspace
EnvironmentFile=-{{INSTALL_DIR}}/.env
# Prepend ~/.local/bin so `claude` (npm-installed there) is on PATH.
Environment=PATH={{LOCAL_BIN}}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart={{NODE_BIN}} {{INSTALL_DIR}}/workspace/bin/cron-runner.js %i
StandardOutput=append:{{INSTALL_DIR}}/logs/%i.log
StandardError=append:{{INSTALL_DIR}}/logs/%i.log
TimeoutStartSec=7200
# Reflections are awaited in-process by cron-runner.js (lib/cron-reflect.js);
# nothing detaches any more. The default KillMode=control-group is correct:
# stop/timeout SIGTERMs the whole cgroup, reaping the Agent SDK's Claude
# Code child too. (Existing installs: regenerate units via gen-units.sh and
# `systemctl daemon-reload`.)
